# iso/config/airootfs/etc/systemd/system/installer.service
[Unit]
Description=installer

[Service]
User=root
ExecStart=startx
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target